# media-gfx/brscan-skey-bin/files/brscan-skey
#!/sbin/openrc-run
# Copyright 2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="brscan-skey daemon"
description="Brother s-key tool daemon"
command=/opt/brother/scanner/brscan-skey/brscan-skey
command_args="${brscan-skey_args}"

depend() {
	after hotplug modules
	need dbus
}

start() {
	ebegin "Starting ${SVCNAME}"
	${command}
	eend $?
}

stop() {
	ebegin "Stopping ${SVCNAME}"
	${command} --terminate
	eend $?
}

reload() {
	ebegin "Reloading ${SVCNAME} settings"
	${command} --refresh
	eend $?
}


